# greenion-agents/build/systemd-service/greenion-agent-server.service
[Unit]
Description=Launches the greenion server agent at computer startup
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/greenion-server

[Install]
WantedBy=graphical.target
